# rails-nginx-unicorn/supervisord.conf
[supervisord]
user=rails
nodaemon=true

[program:unicorn]
command=/bin/bash -c "sudo service unicorn start"

[program:nginx]
command=/bin/bash -c "sudo service nginx restart"
